# vdsm/vdsm-logrotate.service
[Unit]
Description=vdsm hourly logrotate

[Service]
Type=oneshot
ExecStart=/usr/lib/vdsm/vdsm-logrotate
User=vdsm
